# kld/ufw/banip.service
[Unit]
Description=banip

[Service]
Type=simple
TimeoutSec=infinity
ExecStart=/bin/bash /usr/sbin/banip

