# vps/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.FoSVb8/crontab installed on Thu Oct 20 22:01:01 2022)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
3 * * * * /opt/bin/banip
7 5 * * * /opt/bin/mkrules
